# .service
[Unit]
Description=<a description of your application>
After=network.target

[Service]
User= #dockerid
WorkingDirectory=<path to your app>
ExecStart= <app start command>
Restart=on-failure

[Install]
WantedBy=multi-user.target